# plugins/lm-mem/systemd/lm-mem-web.service
[Unit]
Description=lm-mem Web UI (端口 7531)
Documentation=https://github.com/laomou/skills/tree/main/plugins/lm-mem
After=network.target lm-mem-backend.service
Requires=lm-mem-backend.service

[Service]
Type=simple
ExecStart=%h/work/code/skills/plugins/lm-mem/.venv/bin/python %h/work/code/skills/plugins/lm-mem/web.py
Restart=on-failure
RestartSec=3
Environment=MEMORY_CHROMA_URL=http://127.0.0.1:8901
Environment=LM_MEM_WEB_HOST=127.0.0.1
Environment=LM_MEM_WEB_PORT=7531

[Install]
WantedBy=default.target